# v3-k8s/scripts/k8s-token-rotate.service
[Unit]
Description=Rotate kubeadm bootstrap token and store in SSM

[Service]
Type=oneshot
ExecStart=/opt/k8s/rotate-join-token.sh
